# docs/systemd/value-steward-tick.timer
[Unit]
Description=Value Steward Local Tick Timer (execution scheduler)

[Timer]
OnCalendar=Mon..Fri 12:30..15:55/5
Timezone=America/New_York
Persistent=true
AccuracySec=1min
Unit=value-steward-tick.service

[Install]
WantedBy=timers.target
