# docs/systemd/value-steward-world.timer
[Unit]
Description=Value Steward World Context Timer

[Timer]
OnCalendar=Mon..Fri 08:55..15:55/5
Timezone=America/New_York
Persistent=true
AccuracySec=1min
Unit=value-steward-world.service

[Install]
WantedBy=timers.target
