# docs/systemd/value-steward-tick.service
[Unit]
Description=Value Steward Local Tick
After=network.target

[Service]
Type=oneshot
User=lukes
WorkingDirectory=/home/lukes/value-steward
EnvironmentFile=-/home/lukes/value-steward/.env
ExecStart=/usr/bin/npm run local:tick:scheduled
