# docs/systemd/value-steward-gpio.service
[Unit]
Description=Value Steward GPIO Bridge
After=network.target

[Service]
Type=simple
User=lukes
WorkingDirectory=/home/lukes/value-steward
EnvironmentFile=-/home/lukes/value-steward/.env
ExecStart=/usr/bin/npm run gpio:daemon
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
